# service/gfx_clock_NCS3186.service
[Unit]
Description=GFX Clock
After=multi-user.target

[Service]
Type=idle
ExecStart=/home/pi/gfx_clock NCS3186

[Install]
WantedBy=multi-user.target