# service/gfx_clock_NCS3148C.service
[Unit]
Description=GFX Clock
After=multi-user.target

[Service]
Type=idle
ExecStart=/home/pi/gfx_clock NCS3148C

[Install]
WantedBy=multi-user.target